# cookbooks/ghost/templates/default/ghost.ini.erb
[program:ghost]
command = ghost run
directory = <%= node['ghost']['homedir'] %>
user = <%= node['ghost']['user'] %>
autostart = true
autorestart = true
redirect_stdout = true
redirect_stderr = true
environment = NODE_ENV="production"
